# src/httpd/build/modules/httpd/files/etc/supervisor/conf.d/httpd_stderr.conf
[program:httpd_stderr]
command = /usr/bin/tail -F /var/log/apache2/stderr.log
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
