# src/httpd/build/modules/httpd/files/etc/supervisor/conf.d/httpd.conf
[program:httpd]
command = /usr/sbin/apache2ctl -D FOREGROUND
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
