# src/apache-2.4/run/modules/run/files/etc/cron.d/certbot
0 0 */30 * * root certbot-auto renew --post-hook "certbot_renew_certificate_post_hook" >> /var/log/cron.log 2>&1
