# build/etc/supervisor/conf.d/apache2.conf
[program:apache2]
command=apache2ctl -D FOREGROUND
autostart=true
autorestart=true