# src/apache-2.4/run/modules/run/files/etc/supervisor/conf.d/certbot.conf
[program:certbot]
command = /usr/local/bin/certbot_obtain_certificate
autostart = true
autorestart = false
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
