# src/build/modules/httpd/files/etc/supervisor/conf.d/httpd.conf
[program:httpd]
command=/usr/sbin/apache2ctl -c "CustomLog /dev/stdout" -D FOREGROUND
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=TERM
