# src/build/modules/httpd/files/etc/supervisor/conf.d/httpd.conf
[program:httpd]
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true
stderr_logfile=/dev/stdout
stdout_logfile=/dev/stdout
